# scripts/emergency-agent.service
[Unit]
Description=Emergency Agent Session (Daytime / Showcase Mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/andrii/lain/agent_project
Environment=PATH=/home/andrii/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=TRIGGER_MODE=emergency
ExecStart=/usr/bin/bash /home/andrii/lain/agent_project/scripts/wake.sh \
  /home/andrii/lain/agent_project/prompts/emergency_goal.txt \
  /home/andrii/lain/agent_project/prompts/persona.txt
# Cap individual sessions at 30 minutes -- emergency mode runs frequently,
# so short focused sessions are better than one long one.
TimeoutStartSec=30min

[Install]
WantedBy=default.target
