# scripts/night-agent.timer
[Unit]
Description=Trigger Night Agent at fixed wake times

[Timer]
# Five wake times: 23:00 probe slot + four post-reset sessions (~1h15m apart).
# 01:10 chosen to give a 10-minute buffer after the 5h usage window resets at 01:00.
# OnCalendar accepts multiple lines -- each fires independently.
OnCalendar=*-*-* 23:00:00
OnCalendar=*-*-* 01:10:00
OnCalendar=*-*-* 02:25:00
OnCalendar=*-*-* 03:40:00
OnCalendar=*-*-* 04:55:00
Persistent=false
# If the machine is asleep/off at the exact moment, do NOT fire late after
# wake -- a 23:00 session firing at 8am defeats the whole point of "night
# only." Persistent=false is intentional here.

[Install]
WantedBy=timers.target
