# scripts/night-agent.service
[Unit]
Description=Night Agent Wake Session
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=andrii
WorkingDirectory=/home/andrii/lain/agent_project
Environment=TRIGGER_MODE=nightly
ExecStart=/usr/bin/bash /home/andrii/lain/agent_project/scripts/wake.sh \
  /home/andrii/lain/agent_project/prompts/goal.txt \
  /home/andrii/lain/agent_project/prompts/persona.txt
# Hard ceiling so a stuck session can never run past the window.
TimeoutStartSec=7h

[Install]
WantedBy=default.target
